# com.sqlproxy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Service identifier -->
    <key>Label</key>
    <string>com.sqlproxy</string>

    <!-- Executable and arguments -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/sql-proxy</string>
        <string>-config</string>
        <string>/usr/local/etc/sql-proxy/config.yaml</string>
    </array>

    <!-- Working directory -->
    <key>WorkingDirectory</key>
    <string>/usr/local/etc/sql-proxy</string>

    <!-- Start at boot -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Keep alive (restart on crash) -->
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <!-- Restart throttling (wait 5 seconds between restarts) -->
    <key>ThrottleInterval</key>
    <integer>5</integer>

    <!-- Environment variables -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>DB_HOST</key>
        <string>your-server.example.com</string>
        <key>DB_USER</key>
        <string>sqlproxy_reader</string>
        <key>DB_NAME</key>
        <string>YourDB</string>
        <!-- For DB_PASSWORD, use keychain or secrets management -->
    </dict>

    <!-- Logging (stdout/stderr to files) -->
    <key>StandardOutPath</key>
    <string>/usr/local/var/log/sql-proxy/sql-proxy.log</string>
    <key>StandardErrorPath</key>
    <string>/usr/local/var/log/sql-proxy/sql-proxy.err</string>

    <!-- Run as specific user (optional - remove for root) -->
    <key>UserName</key>
    <string>_sqlproxy</string>
    <key>GroupName</key>
    <string>_sqlproxy</string>

    <!-- Resource limits -->
    <key>SoftResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>65536</integer>
    </dict>
    <key>HardResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>65536</integer>
    </dict>
</dict>
</plist>
