# sql-proxy.service
[Unit]
Description=SQL Proxy Service
Documentation=https://github.com/your-org/sql-proxy
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=sqlproxy
Group=sqlproxy

# Working directory and executable
WorkingDirectory=/opt/sql-proxy
ExecStart=/opt/sql-proxy/sql-proxy -config /opt/sql-proxy/config.yaml

# Restart policy
Restart=on-failure
RestartSec=5s

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ReadWritePaths=/var/log/sql-proxy

# Resource limits
LimitNOFILE=65536

# Environment
Environment=DB_HOST=your-server.example.com
Environment=DB_USER=sqlproxy_reader
Environment=DB_NAME=YourDB
# For sensitive values, use EnvironmentFile instead:
# EnvironmentFile=/opt/sql-proxy/.env

# Logging (stdout goes to journald)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sql-proxy

[Install]
WantedBy=multi-user.target
